# tests/ci/kafka-service/supervisor_conf/zookeeper.conf
[program:zookeeper]
command=/opt/kafka/bin/zookeeper-server-start.sh /opt/kafka/config/zookeeper.properties
autostart=true
autorestart=true